# src/test/tools/supervisord.conf
[supervisord]
nodaemon=true

[program:thanhtoan]
command=java -jar /data/thanhtoan/app.jar
