# src/tools/supervisord.conf
[supervisord]
nodaemon=true

[program:thanhtoan]
command=java -Xmx512m -Xms128m -XX:MaxRAM=1536m -XX:+UseSerialGC -jar /data/thanhtoan/app.jar
